# libretranslate/rootfs/etc/service.d/libretranslate/run
#!/usr/bin/with-contenv bashio
set -e

log() {
  if [ "${bashio:-false}" != false ]; then
    bashio::log.info "$1"
  else
    echo "$1"
  fi
}

log "LibreTranslate add-on starting"

# Setze Variablen
HOST="${HOST:-0.0.0.0}"
PORT="${PORT:-5000}"
LANGUAGES="${LANGUAGES:-en,es,de,fr,it,zh}"
DEBUG="${DEBUG:-false}"

MODELS_DIR="/share/libretranslate/models"
mkdir -p "$MODELS_DIR"

# Prüfen Modelle
MISSING=false
IFS=',' read -ra LANG_LIST <<< "$LANGUAGES"
for lang in "${LANG_LIST[@]}"; do
  lang_trimmed="$(echo "$lang" | xargs)"
  if [ -z "$lang_trimmed" ]; then
    continue
  fi
  if [ ! -d "$MODELS_DIR/$lang_trimmed" ]; then
    log "Missing model for: $lang_trimmed"
    MISSING=true
  else
    log "Model exists for: $lang_trimmed"
  fi
done

if [ "$MISSING" = true ]; then
  log "Downloading missing models..."
  if [ "$DEBUG" = true ] || [ "$DEBUG" = "true" ]; then
    libretranslate --update-models --models-dir "$MODELS_DIR" --debug
  else
    libretranslate --update-models --models-dir "$MODELS_DIR"
  fi
else
  log "All configured models present"
fi

log "Starting LibreTranslate on ${HOST}:${PORT}"
exec libretranslate --host "${HOST}" --port "${PORT}" --models-dir "${MODELS_DIR}" --load-only "${LANGUAGES}" $([ "$DEBUG" = true ] && echo "--debug")
